# deploy/dybervpn.service
[Unit]
Description=DyberVPN Post-Quantum VPN Daemon
Documentation=https://github.com/dyberinc/dybervpn
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/dybervpn up -c /etc/dybervpn/config.toml -f
ExecStop=/usr/local/bin/dybervpn down dvpn0
Restart=on-failure
RestartSec=5

# Security hardening
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_BIND_SERVICE CAP_NET_RAW
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_BIND_SERVICE CAP_NET_RAW
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
ReadWritePaths=/run

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dybervpn

[Install]
WantedBy=multi-user.target
